# deploy/my-aibox.service
# MyAIBOX systemd user service
# Do not install this file directly. Run `my-aibox install`.

[Unit]
Description=MyAIBOX Service
After=network-online.target

[Service]
Type=simple
WorkingDirectory={project_dir}
ExecStart={project_dir}/.venv/bin/python app.py
EnvironmentFile={project_dir}/.env
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
